# rules.d/10-xppen.rules
#
#
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="28bd", ATTRS{idProduct}=="0075", RUN+="/bin/sh -e /etc/udev/rebind.sh %p"